# lib/dockerfiles/graphite-web/files/supervisord.conf
[supervisord]
nodaemon = true

[program:postgres]
priority = 100
user = postgres
command = /usr/lib/postgresql/9.4/bin/postgres -D /var/lib/postgresql/9.4/main -c config_file=/etc/postgresql/9.4/main/postgresql.conf
autostart = true
autorestart = true
startsecs = 10

[eventlistener:supervisord_event_listener]
priority = 200
command = /root/supervisord_event_listener.py
autostart = true
autorestart = false
events = PROCESS_STATE_RUNNING

[program:uwsgi]
priority = 300
command = /usr/bin/uwsgi --emperor /etc/uwsgi/apps-enabled
autostart = false
autorestart = true
startsecs = 10
stopsignal = INT

[program:nginx]
priority = 400
command = /usr/sbin/nginx
autostart = false
autorestart = true
startsecs = 10
stdout_logfile = NONE
