# elasticsearch/files/supervisord.conf
[supervisord]
nodaemon = true

[program:elasticsearch]
user = elasticsearch
command = /usr/share/elasticsearch/bin/elasticsearch --default.path.data=/var/lib/elasticsearch
autostart = true
autorestart = true
startsecs = 10
environment = ES_HEAP_SIZE="4096m"
